# openshift/drupal-repository/files/etc/s6-overlay/s6-rc.d/webhook/run
#!/command/with-contenv sh
/usr/local/bin/webhook -verbose -hooks=/etc/webhook/hooks.json -template -hotreload -debug
