# openshift/drupal-repository/files/etc/s6-overlay/s6-rc.d/package-queue/run
#!/command/with-contenv sh
cd $PROJECT_DIR
make
exec php console.php queue:consume package-queue --item-limit 10;
